# dynserver.service
[Unit]
Description=Dynamic Gallery Server
Documentation=https://github.com/your-org/dynserver
After=network.target network-online.target
Requires=network-online.target

[Service]
Type=exec
User=dynserver
Group=dynserver
WorkingDirectory=/usr/share/dynserver
ExecStart=/usr/bin/dynserver --config /etc/dynserver/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
KillMode=mixed

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

[Install]
WantedBy=multi-user.target
