# tenrankai.service
[Unit]
Description=Dynamic Gallery Server
Documentation=https://github.com/theatrus/tenrankai
After=network.target network-online.target
Requires=network-online.target

[Service]
Type=exec
User=tenrankai
Group=tenrankai
WorkingDirectory=/usr/share/tenrankai
ExecStart=/usr/bin/tenrankai --config /etc/tenrankai/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
KillMode=mixed

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

[Install]
WantedBy=multi-user.target
