# debian/postinst
#!/bin/bash

set -e

# Source debconf library
. /usr/share/debconf/confmodule

case "$1" in
    configure)
        # Create tenrankai system user and group
        if ! getent group tenrankai >/dev/null 2>&1; then
            addgroup --system tenrankai
        fi
        
        if ! getent passwd tenrankai >/dev/null 2>&1; then
            adduser --system --group --home /var/lib/tenrankai \
                    --shell /bin/false --gecos "DynServer daemon" \
                    --disabled-password tenrankai
        fi
        
        # Set ownership and permissions for directories
        chown -R tenrankai:tenrankai /var/lib/tenrankai
        chown -R tenrankai:tenrankai /var/log/tenrankai
        chown -R root:tenrankai /etc/tenrankai
        
        # Set correct permissions
        chmod 750 /etc/tenrankai
        chmod 750 /var/lib/tenrankai
        chmod 755 /var/log/tenrankai
        
        # Create example config if no config exists
        if [ ! -f /etc/tenrankai/config.toml ]; then
            echo "No configuration file found. Please copy and customize:"
            echo "  sudo cp /etc/tenrankai/config.toml.example /etc/tenrankai/config.toml"
            echo "  sudo chown root:tenrankai /etc/tenrankai/config.toml"
            echo "  sudo chmod 640 /etc/tenrankai/config.toml"
            echo ""
            echo "Then edit /etc/tenrankai/config.toml to match your setup."
        fi
        
        # Reload systemd and enable service
        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload
            # Don't auto-start until configured
            echo "Service installed but not started. Configure first, then:"
            echo "  sudo systemctl enable tenrankai"
            echo "  sudo systemctl start tenrankai"
        fi
        ;;
esac

#DEBHELPER#

exit 0